# armeiroMaster/com.visualstudio.code.tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
<key>Label</key>
<string>com.visualstudio.code.tunnel</string>
<key>LimitLoadToSessionType</key>
<string>Aqua</string>
<key>ProgramArguments</key>
<array>
<string>/private/var/folders/cv/zm431vg57gs62jw8c4d60q6r0000gn/T/AppTranslocation/60E84047-8945-4B3B-9C77-05D5C17B187B/d/Visual Studio Code.app/Contents/Resources/app/bin/code-tunnel</string>
<string>--verbose</string><string>--cli-data-dir</string><string>/Users/saranardelli/.vscode/cli</string><string>tunnel</string><string>service</string><string>internal-run</string>
</array>
<key>KeepAlive</key>
<true/>
<key>StandardErrorPath</key>
<string>/Users/saranardelli/.vscode/cli/tunnel-service.log</string>
<key>StandardOutPath</key>
<string>/Users/saranardelli/.vscode/cli/tunnel-service.log</string>
</dict>
</plist>